# deploy/example-query-mailer-replay.service
[Unit]
Description=Replay failed query-mailer notifications
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__SERVICE_USER__
Group=__SERVICE_GROUP__
WorkingDirectory=__APP_ROOT__/app
EnvironmentFile=__ENV_FILE__
UMask=0077
ExecStart=__VENV_DIR__/bin/python -m query_mailer.replay_failed --limit 100
